# root/etc/services.d/openldap/run
#!/usr/bin/with-contenv bash

# Need to figure out how to run slapd ... I think part of this will be figuring out where the config file will live
